# files/system/usr/lib/systemd/system/greetd.service.d/override.conf
[Service]
TTYPath=/dev/tty2
TTYReset=true
TTYVHangup=true
TTYVTDisallocate=true
