# files/system/etc/udev/rules.d/90-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add", RUN+="/bin/chmod a+w /sys/class/backlight/%k/brightness"
